# deploy/gcp/vendorverdict-monitor.timer
[Unit]
Description=Run VendorVerdict health checks every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true
Unit=vendorverdict-monitor.service

[Install]
WantedBy=timers.target
